# docs/ops/systemd-user/codex-feishu-bridge-healthcheck.service
[Unit]
Description=Healthcheck for codex-feishu-bridge (user)

[Service]
Type=oneshot
ExecStart=__HEALTHCHECK_SCRIPT__
